# chef/cookbooks/baseconfig/files/default/rc.local
#!/bin/sh
sleep 10
uwsgi --ini /home/ubuntu/project/uwsgi.ini  --daemonize /home/ubuntu/mysite.log
